# hisense_mqtt_bridge/rootfs/etc/services.d/hisense-mqtt/run
#!/usr/bin/with-contenv bashio

exec python3 /app/hisense_mqtt_bridge.py
